# nginx-uwsgi/supervisor-app.conf
[program:app-uwsgi]
command = /usr/bin/uwsgi /var/www/bottle/app.ini

[program:nginx-app]
command = /usr/sbin/nginx
